# deploy/pm_shadow.service
[Unit]
Description=Polymarket shadow copy-trading (no orders)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pm-shadow/Son-of-weather
ExecStart=/opt/pm-shadow/Son-of-weather/.venv/bin/python /opt/pm-shadow/Son-of-weather/run_shadow.py --config /opt/pm-shadow/Son-of-weather/config.yaml --out-dir /opt/pm-shadow/Son-of-weather/reports --state-db /opt/pm-shadow/Son-of-weather/shadow_state.db --days 7
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# If your API requests must go through Clash/HTTP proxy, uncomment:
# Environment=http_proxy=http://127.0.0.1:7897
# Environment=https_proxy=http://127.0.0.1:7897
# Environment=HTTP_PROXY=http://127.0.0.1:7897
# Environment=HTTPS_PROXY=http://127.0.0.1:7897
# Environment=no_proxy=localhost,127.0.0.1
# Environment=NO_PROXY=localhost,127.0.0.1

# hardening (optional)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pm-shadow/Son-of-weather

[Install]
WantedBy=multi-user.target

